# rework/core/deploy/supervisor/wistar_supervisor_test.conf
[program:wistar_server_test]
command=uwsgi --ini /opt/test-projects/wistar-server/deploy/uwsgi/wistar_uwsgi_test.ini
environment=DJANGO_SETTINGS_MODULE=wistar.settings.test
stderr_logfile = /var/log/uwsgi/wistar_server_test_stderr.log
stdout_logfile = /var/log/uwsgi/wistar_server_test_stdout.log
