# rework/core/deploy/supervisor/wistar_supervisor_production.conf
[program:wistar_server]
command=uwsgi --ini /opt/projects/wistar-server/deploy/uwsgi/wistar_uwsgi_production.ini
environment=DJANGO_SETTINGS_MODULE=wistar.settings.production
stderr_logfile = /var/log/uwsgi/wistar_server_stderr.log
stdout_logfile = /var/log/uwsgi/wistar_server_stdout.log
